# deploy/lagoon.service
[Unit]
Description=lagoon notes web app
Documentation=https://github.com/deepwa7er/lagoon
# Order after tailscaled so the Tailscale interface/IP exists before we bind.
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=lagoon
Group=lagoon
ExecStart=/usr/local/bin/lagoon-server --config /etc/lagoon/config.toml
Restart=on-failure
RestartSec=2

# Hardening — the service serves static files, reads the embedding model, and
# reads/writes only its own SQLite database under /var/lib/lagoon.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
# The store lives here; everything else is read-only under ProtectSystem=strict.
ReadWritePaths=/var/lib/lagoon

[Install]
WantedBy=multi-user.target
